# etc/cron.d/eeepc-tweaks

@reboot		root	[ -e "/usr/sbin/eeepc-tweaks" ] && /usr/sbin/eeepc-tweaks
